# systemd/cert-renewer.target
[Unit]
Description=Target to restart all parts of the Smallstep Certificate Renewal
