# systemd/cert-renewer.target
[Unit]
Description=Target to restart all certificate renewers
